# landos/scripts/com.basemodhomes.landos-pipeline.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.basemodhomes.landos-pipeline</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/bennett2026/Desktop/Kingdom LandOS/BaseMod LandOS/landos/scripts/daily_pipeline.sh</string>
    </array>

    <!-- Run daily at 6:00 AM local time -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/Users/bennett2026/Desktop/Kingdom LandOS/BaseMod LandOS</string>

    <!-- Stdout/stderr for launchd-level logging -->
    <key>StandardOutPath</key>
    <string>/Users/bennett2026/Desktop/Kingdom LandOS/BaseMod LandOS/landos/logs/launchd_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/bennett2026/Desktop/Kingdom LandOS/BaseMod LandOS/landos/logs/launchd_stderr.log</string>

</dict>
</plist>
